# config/u2f/u2f-session.service
[Unit]
Description=Lock session when YubiKey is removed (unless Super is held)
After=graphical-session.target
BindsTo=graphical-session.target

[Service]
Type=simple
Environment=PYTHONUNBUFFERED=1
ExecStart=%h/.local/lib/u2f-session/u2f_session_monitor.py
Restart=on-failure
RestartSec=2s

[Install]
WantedBy=default.target
